# dot_config/systemd/user/random-feh.timer
[Unit]
Description=Set a random wallpaper every 5 minutes

[Timer]
OnBootSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=default.target
